# systemd/polymarket-copybot-discord-monitor.service
# Polymarket Copybot Discord monitor
# Read-only: polls loopback status API, posts alerts/digests, and exposes /status /wallets /gaps only.
[Unit]
Description=Polymarket Copybot Discord Monitor
After=network-online.target polymarket-copybot-status-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/flip/projects/polymarket-copybot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/root/flip/secrets/copybot_discord.env
ExecStart=/usr/local/lib/hermes-agent/venv/bin/python -m polymarket_bot.discord_monitor
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
