# systemd/polymarket-copybot-status-api.service
# Polymarket Copybot read-only status API
# Loopback-only; GET endpoints only; no keys, no orders, no controls.
[Unit]
Description=Polymarket Copybot Read-Only Status API
After=network-online.target polymarket-copybot-book-archive.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/flip/projects/polymarket-copybot
Environment=PYTHONUNBUFFERED=1
Environment=POLYMARKET_ARCHIVE_CONFIG=/root/flip/projects/polymarket-copybot/archive_config.json
Environment=POLYMARKET_STATUS_HOST=0.0.0.0
Environment=POLYMARKET_STATUS_PORT=8710
ExecStart=/usr/bin/python3 -m polymarket_bot.status_api
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
