# systemd/polymarket-copybot-paper-follower.service
# Polymarket Copybot paper follower daemon
# PAPER ONLY. Public data only. No keys, no orders, no live execution paths.
[Unit]
Description=Polymarket Copybot Paper Follower
After=network-online.target polymarket-copybot-book-archive.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/flip/projects/polymarket-copybot
Environment=PYTHONUNBUFFERED=1
Environment=POLYMARKET_ARCHIVE_CONFIG=/root/flip/projects/polymarket-copybot/archive_config.json
ExecStart=/usr/bin/python3 -m polymarket_bot.paper_follower
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/root/flip/projects/polymarket-copybot/runs

[Install]
WantedBy=multi-user.target
